# etc/profile.d/profile.sh
#!/usr/bin/env bash

[ -s "$HOME/etc/profile.d/aliases.sh" ] && . "$HOME/etc/profile.d/aliases.sh"
PATH="$HOME/.yarn/bin:$HOME/anaconda3/bin:$PATH"
CDPATH=.:~:~/github:~/github/kaluzki:~/gitlab